# systemd/pisugar-button.service
[Unit]
Description=PiSugar S Button Daemon (GPIO3)
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 %h/lotus-companion/scripts/pisugar_button.py
Environment=GPIOZERO_PIN_FACTORY=lgpio
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
